# ops/systemd/wheatandstone-fulfillment-automation.service
[Unit]
Description=Wheat & Stone fulfillment automation scheduler
After=network.target wheatandstone-app.service
Wants=wheatandstone-app.service

[Service]
Type=oneshot
WorkingDirectory=/var/www/WheatAndStone/ws-app
EnvironmentFile=/etc/wheatandstone/wheatandstone-web.env
ExecStart=/usr/bin/env bash -lc 'pnpm automation:fulfillment:run'
